# playgrounds/caddyserver/image/caddy/caddy-trust.service
[Unit]
Description=Install Caddy local CA into system trust store (after Caddy starts)
After=caddy.service
PartOf=caddy.service

[Service]
Type=oneshot
ExecStart=/usr/bin/caddy trust

[Install]
WantedBy=caddy.service
