# images/cri/image/_rootfs/etc/systemd/system/containerd.service.d/10-container-runtime.conf
[Unit]
Wants=container-runtime.target
Before=container-runtime.target
