# images/cri/image/_rootfs/etc/systemd/system/container-runtime.target
[Unit]
Description=Generic a container runtime is up target
